# djangosite.conf
[program:djangosite]
directory=/app/djangosite
command=/app/venv/bin/uwsgi --ini /app/uwsgi.ini
stdout_logfile=/app/logs/supervisor.log